# app/misc/oxomium.socket
[Unit]
Description=oxomium socket

[Socket]
ListenStream=/run/oxomium.sock
SocketUser=www-data
SocketMode=600

[Install]
WantedBy=sockets.target
